# etc/yasa/etc/10-usblock.rules
# Do not edit this file. This fiel was automatically generated by /usr/sbin/usblock_rule_generator
# Based on article from <Daniel KopeÄek> in GitHub.
#SUBSYSTEM=="usb", TEST=="authorized", RUN+="/etc/usblock/rule_generator pipe  %E{ID_VENDOR_ID} %E{ID_MODEL_ID} %E{ID_SERIAL_SHORT}"
# Send usb device information to syslog

ACTION!="add", GOTO="oakusb_rules_end"
SUBSYSTEM!="usb", GOTO="oakusb_rules_end"
TEST!="authorized", GOTO="oakusb_rules_end"

# log usb device first time connect and authorized device
RUN+="/usr/bin/logger -p local0.info 'Device:%E{ID_MODEL} VID:%E{ID_VENDOR_ID} PID:%E{ID_MODEL_ID} Serial:%E{ID_SERIAL_SHORT}'"

## make hubs deauthorized all devices by default
TEST=="authorized_default", ATTR{authorized_default}="0", GOTO="oakusb_rules_end"

ENV{valid}=="1", ENV{valid}="", ATTR{authorized}="1", GOTO="oakusb_rules_end"

LABEL="oakusb_rules_end"
